# deploy/systemd/pc-dashboard.service
[Unit]
Description=PC Backend (Django runserver, auto-reload on code change)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Desktop/pc-dashboard
ExecStart=%h/miniconda3/bin/python manage.py runserver 0.0.0.0:8000
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
